# Linuxforadxl/build/linux/rootfs/packages-repo/targetroot/var/lib/dpkg/info/busybox-inetd.list
/.
/etc
/etc/inetd.conf
/etc/init.d
/etc/init.d/inetd.busybox
