# Linuxforadxl/build/linux/rootfs/packages-repo/targetroot/var/lib/dpkg/info/busybox-mdev.list
/.
/etc
/etc/mdev.conf
/etc/init.d
/etc/init.d/mdev
